# src/build/etc/supervisor/conf.d/apache.conf
[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true