# build/etc/supervisor/conf.d/php.conf
[program:php]
command=/phpfarm/inst/php-5.4.33/sbin/php-fpm -F
autostart=true
autorestart=true
